# data/passim.service.in
[Unit]
Description=Local Caching Server
Documentation=https://github.com/hughsie/passim
After=avahi-daemon.service
Before=display-manager.service
Wants=avahi-daemon.service

[Service]
Type=dbus
TimeoutSec=180
BusName=org.freedesktop.Passim
ExecStart=@libexecdir@/passimd
User=passim
DevicePolicy=closed
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateMounts=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
# allow-list, see `systemd-analyze syscall-filter` for what each group adds
SystemCallFilter=@basic-io @file-system @io-event @ipc @network-io @process @signal madvise
SystemCallErrorNumber=SIGSYS
SystemCallArchitectures=native
StateDirectory=passim passim/data
LogsDirectory=passim
