# data/passim.service.in
[Unit]
Description=A local caching server
Documentation=https://github.com/hughsie/passim
After=dbus.service
Before=display-manager.service

[Service]
Type=dbus
TimeoutSec=180
BusName=org.freedesktop.Passim
ExecStart=@libexecdir@/passimd
DevicePolicy=closed
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateMounts=true
PrivateTmp=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=full
ReadWritePaths=-@localstatedir@/lib/passim/data
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallFilter=~@mount
