# data/passim.service.in
[Unit]
Description=Local Caching Server
Documentation=https://github.com/hughsie/passim
After=avahi-daemon.service
Before=display-manager.service
Wants=avahi-daemon.service

[Service]
Type=dbus
TimeoutSec=180
BusName=org.freedesktop.Passim
ExecStart=@libexecdir@/passimd
User=passim
DevicePolicy=closed
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateMounts=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
StateDirectory=passim passim/data
LogsDirectory=passim
